# packaging/energybuildai-dashboard.service
[Unit]
Description=EnergyBuildAI Qt Dashboard
After=graphical.target network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment=QT_QPA_PLATFORM=wayland
Environment=QT_QUICK_CONTROLS_STYLE=Fusion
Environment=QT_WAYLAND_DISABLE_WINDOWDECORATION=1
ExecStart=/usr/bin/energybuildai-dashboard
Restart=always
RestartSec=2

[Install]
WantedBy=graphical.target
